# deploy/report-backend.service
[Unit]
Description=RePort Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=report
Group=report
WorkingDirectory=/home/report/RePort/backend
Environment="PATH=/home/report/venv/bin"
ExecStart=/home/report/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/report/backend.log
StandardError=append:/var/log/report/backend-error.log

[Install]
WantedBy=multi-user.target
